# system/speedtest.service
[Unit]
Description=Speedtest Service

[Service]
Type=oneshot
ExecStart=/srv/tachyon/scripts/speedtest.sh
